# webapp/supervisord.conf
# =============================================================================
# Supervisor Configuration for MetaSoccer ID Web Application
# =============================================================================
# This configuration manages multiple processes in the Docker container:
# - nginx: Reverse proxy server (priority 10 - starts first)
# - nextjs: Next.js application server (priority 20 - starts after nginx)
# =============================================================================

# -----------------------------------------------------------------------------
# MAIN SUPERVISOR DAEMON CONFIGURATION
# -----------------------------------------------------------------------------
[supervisord]
# Run supervisor in foreground (required for Docker containers)
nodaemon=true

# Run supervisor as root user (required for managing other users)
user=root

# Logging configuration
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info

# Process management
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

# Prevent supervisor from clearing environment variables
strip_ansi=false

# -----------------------------------------------------------------------------
# UNIX HTTP SERVER CONFIGURATION
# -----------------------------------------------------------------------------
# Enables supervisorctl to communicate with supervisord
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
chown=root:root

# -----------------------------------------------------------------------------
# SUPERVISORCTL CONFIGURATION
# -----------------------------------------------------------------------------
# Configuration for the supervisorctl client
[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# -----------------------------------------------------------------------------
# RPC INTERFACE CONFIGURATION
# -----------------------------------------------------------------------------
# Required for supervisorctl to work
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# -----------------------------------------------------------------------------
# PROGRAM CONFIGURATIONS
# -----------------------------------------------------------------------------

# =============================================================================
# NGINX REVERSE PROXY SERVER
# =============================================================================
# Handles HTTP requests and proxies API calls to external server
[program:nginx]
# Command to run nginx in foreground mode
command=nginx -g 'daemon off;'

# Process management
priority=10
autostart=true
autorestart=true
startsecs=10
startretries=3

# User and directory settings
user=nginx
directory=/etc/nginx

# Logging configuration
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5

# Environment variables for nginx
environment=PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Signal handling
stopsignal=QUIT
stopwaitsecs=10

# =============================================================================
# NEXT.JS APPLICATION SERVER
# =============================================================================
# Serves the Next.js web application
[program:nextjs]
# Command to run the Next.js server
command=node server.js

# Process management
priority=20
autostart=true
autorestart=true
startsecs=10
startretries=3

# User and directory settings
user=nextjs
directory=/app

# Logging configuration
stdout_logfile=/var/log/supervisor/nextjs.log
stderr_logfile=/var/log/supervisor/nextjs.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5

# Environment variables for Next.js
environment=NODE_ENV=production,
           PORT=3000,
           HOSTNAME="0.0.0.0",
           PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Signal handling
stopsignal=TERM
stopwaitsecs=30

# -----------------------------------------------------------------------------
# EVENT LISTENERS (Optional)
# -----------------------------------------------------------------------------
# Uncomment the following section to enable email notifications on failures
# [eventlistener:crashmail]
# command=crashmail -a -m user@example.com
# events=PROCESS_STATE_FATAL
# priority=1

# -----------------------------------------------------------------------------
# PROCESS GROUPS (Optional)
# -----------------------------------------------------------------------------
# Uncomment to group related processes together
# [group:webapp]
# programs=nginx,nextjs
# priority=10 